# casino-robot.service
[Unit]
Description=Casino Robot Web Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/var/www/casino-robot
ExecStart=/usr/bin/python3 serve.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
